# system/BaNG-Trigger.service
[Unit]
Description=Start Backup-Trigger script

[Service]
Type=oneshot
ExecStart=/opt/BaNG/BaNG-Trigger
